# samarassr.conf
[program:samarassr]
directory=/var/www/samara.medin.cloud
command=node server-dist/index.js
numprocs=1
autostart=true
autorestart=true
user=www-data